# reference/scripts/com.thomaslee.twitter-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- REFERENCE ONLY. The live LaunchAgent is generated by scripts/schedule-launchd.sh
     from the actual repo location. This file shows what that generated plist looks like
     (paths assume the project lives at ~/automations/twitter_summary). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thomaslee.twitter-digest</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/thomaslee/automations/twitter_summary/scripts/run-digest.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>0</integer>
    </dict>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/thomaslee/automations/twitter_summary/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/thomaslee/automations/twitter_summary/logs/launchd-stderr.log</string>
</dict>
</plist>
